# mongo/supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/opt/mongodb-linux-x86_64-3.0.4/bin/mongod --config /opt/mongodb-linux-x86_64-3.0.4/mongo.conf
